# rootfs/etc/services.d/syslog/run
#!/bin/sh

exec syslogd -n -a /var/spool/postfix/dev/log
